# systemd/kronos-daily-status.timer
[Unit]
Description=Kronos Agent OS Daily Status Report

[Timer]
OnCalendar=*-*-* 23:00:00
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
