# systemd/kronos-health.service
[Unit]
Description=Kronos Agent OS Health Check
After=kaos.service

[Service]
Type=oneshot
ExecStart=/opt/kaos/app/scripts/health-check.sh --alert
EnvironmentFile=/opt/kaos/app/.env

StandardOutput=journal+console
StandardError=journal+console
SyslogIdentifier=kronos-health

# Alert on failure
OnFailure=kronos-alert@%n.service
