# systemd/kronos-daily-status.service
[Unit]
Description=Kronos Agent OS Daily Status Report
After=kaos.service

[Service]
Type=oneshot
ExecStart=/opt/kaos/app/scripts/daily-status.sh
EnvironmentFile=/opt/kaos/app/.env

StandardOutput=journal+console
StandardError=journal+console
SyslogIdentifier=kronos-daily-status
