# systemd/kronos-ii.service
[Unit]
Description=Kronos II — LangGraph AI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kronos
WorkingDirectory=/opt/kronos-ii/app
ExecStart=/opt/kronos-ii/app/.venv/bin/python -m kronos
Restart=on-failure
RestartSec=10

# Kill ALL child processes (MCP servers, threads) on stop/restart
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=30

EnvironmentFile=/opt/kronos-ii/app/.env

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kronos-ii

[Install]
WantedBy=multi-user.target
