# systemd/kronos-health.timer
[Unit]
Description=Kronos Agent OS Health Check (every 15 min)

[Timer]
OnCalendar=*:0/15
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
