# systemd/kronos-backup.service
[Unit]
Description=Kronos Agent OS Workspace Backup
After=kaos.service

[Service]
Type=oneshot
ExecStart=/opt/kaos/app/scripts/workspace-backup.sh
EnvironmentFile=/opt/kaos/app/.env

StandardOutput=journal+console
StandardError=journal+console
SyslogIdentifier=kronos-backup
