# systemd/lacuna.service
[Unit]
Description=Lacuna — Creative Director Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kronos
WorkingDirectory=/opt/kronos-ii/app
ExecStart=/opt/kronos-ii/app/.venv/bin/python -m kronos
Restart=on-failure
RestartSec=10

EnvironmentFile=/opt/kronos-ii/app/.env.lacuna
Environment=WEBHOOK_PORT=8790 DASHBOARD_PORT=18790 WORKSPACE_PATH=

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lacuna

[Install]
WantedBy=multi-user.target
